# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
# Redirect main supervisor logs to null to prevent double-logging in Docker logs
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[program:octane]
process_name=%(program_name)s
# Use octane:frankenphp directly for clarity.
# --workers=auto sets workers based on CPU cores.
# --max-requests=500 prevents memory leaks by restarting workers periodically.
command=php /app/artisan octane:frankenphp --host=0.0.0.0 --port=8000 --workers=2 --max-requests=500
autostart=true
autorestart=true
# Give Octane time to shut down gracefully (3600s is Laravel's recommended default)
stopwaitsecs=3600
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:worker]
process_name=%(program_name)s_%(process_num)02d
command=php /app/artisan queue:work --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=2
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
